# example/crontab
* * * * * [Name="HTTP Test" Executor=http] [URL="https://httpbin.org/get" Method=get]
